# .../self/etc/systemd/system/ellipsis@.timer
[Unit]
Description=ellipsis (%I) ...

[Timer]
OnBootSec=15min
OnUnitInactiveSec=8h

[Install]
WantedBy=multi-user.target
